# demo-hw-amd-vpr/mixins/x86/etc/network/if-pre-up.d/xenbr0
#!/bin/sh

if [ "$IFACE" = xenbr0 ]; then
    brctl addbr xenbr0 || /bin/true
    brctl addif xenbr0 eth0 || /bin/true
    ifconfig eth0 up
fi
